# TMR2026/systemd/carrito_tmr.service
[Unit]
Description=Carrito TMR 2026 — Sistema Autonomo
After=graphical.target bluetooth.target
Wants=bluetooth.target

[Service]
User=angel01
Group=angel01
WorkingDirectory=/home/angel01/Carrito/TMR2026
Environment="PYTHONUNBUFFERED=1"
Environment="SDL_VIDEODRIVER=dummy"
Environment="SDL_AUDIODRIVER=dummy"
Environment="DISPLAY=:0"
Environment="XAUTHORITY=/home/angel01/.Xauthority"
ExecStart=/usr/bin/python3 /home/angel01/Carrito/TMR2026/main.py --display
Restart=on-failure
RestartSec=5s
TimeoutStopSec=10s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=carrito_tmr

[Install]
WantedBy=graphical.target
